# etc/systemd/user/update-monitoring.service
[Unit]
Description=Update Monitoring

[Service]
# Type=oneshot
ExecStart=/usr/local/bin/update-monitoring

[Install]
WantedBy=default.target
