# etc/systemd/system/update-monitoring.path
[Unit]
Description=Update Monitoring

[Path]
PathExists=/var/lib/pacman/db.lck
Unit=update-monitoring.service

[Install]
WantedBy=default.target
